# scripts/line-bot-vps.service
[Unit]
Description=LINE Bot VPS Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/AI-Vtuber-Project
Environment="PATH=/home/ubuntu/AI-Vtuber-Project/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=/home/ubuntu/AI-Vtuber-Project/.env
ExecStart=/home/ubuntu/AI-Vtuber-Project/venv/bin/python -m uvicorn src.line_bot_vps.webhook_server_vps:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
